# cron/crontab
@hourly /usr/local/bin/cron-update-email
@hourly /usr/local/bin/cron-update-rss
